# deploy/conf.d/salasvirtuales_dev.ini
[program:salasvirtuales-web]
command = /apps/salasvirtuales/deploy/scripts/web_dev.sh
user = nginx
autostart = true
autorestart = true
redirect_stderr = true
stdout_logfile_maxbytes = 1MB
stdout_logfile = /apps/salasvirtuales/deploy/logs/web.log

[group:salasvirtuales]
programs = salasvirtuales-web
